# srcs/mysql/supervisord.conf
[supervisord]
nodaemon=true
[program:mysqld_safe]
command = sh -c "/usr/bin/mysqld_safe --user=mysql --datadir='/var/lib/mysql' && kill -s SIGTERM $(cat supervisord.pid)"
[program:telegraf]
command = sh -c "telegraf && kill -s SIGTERM $(cat supervisord.pid)"
